# consul/consul.upstart.conf
# /etc/inid/consul.conf

start on filesystem and static-network-up

exec consul agent -ui -config-dir /etc/consul
